# ledger-agent/deploy/ledger-mcp.service
[Unit]
Description=Ledger MCP server (read tier)
After=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ledger-agent
Environment=LEDGER_DB=/home/pi/financeapp/finance.db
Environment=LEDGER_MCP_PORT=8091
# Environment=LEDGER_SPLIT_MODE=payer_share   <- uncomment ONLY after verifying (see README)
ExecStart=/usr/bin/python3 /home/pi/ledger-agent/ledger_mcp.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
