# pifinance/deploy/pifinance-sync.service
# One-shot SimpleFIN sync — /etc/systemd/system/pifinance-sync.service
# Triggered by pifinance-sync.timer; not enabled directly.
[Unit]
Description=Pi Finance SimpleFIN daily sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/pifinance
EnvironmentFile=/home/pi/pifinance/.env
ExecStart=/home/pi/pifinance/venv/bin/python /home/pi/pifinance/simplefin_sync.py
