# pifinance/deploy/pifinance-sync.timer
# Daily sync schedule — /etc/systemd/system/pifinance-sync.timer
[Unit]
Description=Run Pi Finance SimpleFIN sync once a day

[Timer]
OnCalendar=*-*-* 06:30:00
# If the Pi was off at 06:30, run when it next boots
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
